# bugzilla.cron
#!/bin/sh
cd /usr/share/bugzilla
./collectstats.pl
./whineatnews.pl
