# bugzilla.cron
#!/bin/sh
cd /usr/share/bugzilla
./collectstats.pl
./whineatnews.pl
echo "*** Bugzilla nighlty flush has done! ***"
